# Service/tshirtservice.conf
[program:tshirtservice]
directory=/home/pi/StartfeldSmartTshirt/Service
command=/bin/bash -E -c ./start.sh
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
user=root