# Server/tshirtserver.conf
[program:tshirtserver]
directory=~/StartfeldSmartTshirt/Server/src
command=/bin/bash -E -c ./start.sh
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true